# system/conf.d/aria2c.conf
[program:websockify]
command=/system/conf.d/aria2c.sh
autorestart=true
