# system/supervisord.conf
[supervisord]
nodaemon=true

[program:aria2c]
command=bash /system/conf.d/aria2c.sh
autorestart=true

[include]
files = /system/conf.d/*.conf
